# templates/supervisor.conf
[program:webapp]
user={{ app_user }}
command=gunicorn -c {{ remote['gunicorn_conf_path'] }} {{ wsgi_application | mandatory }}
directory={{ project_root_absolute_path }}
