# templates/supervisor.conf
[program:webapp]
user=nobody
command=gunicorn -c {{ remote['gunicorn_conf_path'] }} psngamesweb.wsgi:application
directory={{ project_root_absolute_path }}
